# ww1111.conf
[program:warn1111-mac16inch]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/python3 /home/ec2-user/warn1111/mac16inch.py
autostart=true
autorestart=true
user=ec2-user
numprocs=1
redirect_stderr=true
stdout_logfile=/home/ec2-user/warn1111/mac16inch.log

[program:warn1111-iphone]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/python3 /home/ec2-user/warn1111/iphone.py
autostart=true
autorestart=true
user=ec2-user
numprocs=1
redirect_stderr=true
stdout_logfile=/home/ec2-user/warn1111/iphone.log